# projects/charging-log-analytics/ops/ranking_snapshot.cron.example
CRON_TZ=Europe/Warsaw
0 3 * * * cd /opt/charging-log-app && /usr/bin/env bash ./run_ranking_snapshot.sh >> /opt/charging-log-app/logs/cron.log 2>&1
